# docker/root/etc/s6-overlay/s6-rc.d/kodi/run
#!/command/execlineb -P

with-contenv

kodi